# thermite
#!/sbin/runscript
# thermite init script
# Written by John Roman
#Distributed under the terms of the GNU General Public License v2
# $Header: $

depend() {
	need localmount
 }

start() {
	ebegin "Starting Thermite"
	start-stop-daemon --background --start --exec /usr/local/bin/thermite.py --make-pidfile --pidfile /var/run/thermite.pid
    	eend $?
 }

stop() {
	start-stop-daemon --stop --exec /usr/local/thermite.py --pidfile /var/run/thermite.pid
 }

